# recipes-connectivity/bluez/files/bluetooth-ttyO1.service
[Unit]
Description=Attach CSR Bluetooth Adapter on ttyO1

[Service]
ExecStart=/usr/sbin/hciattach -n ttyO1 csr

[Install]
WantedBy=default.target
