# recipes-core/busybox/busybox/udhcpc@eth0.service
[Unit]
Description=Connection service (udhcpc)

[Service]
ExecStart=/sbin/udhcpc -ieth0
RemainAfterExit=yes

[Install]
WantedBy=network.target
